# meta-inventec/meta-transformers-nuv/recipes-phosphor/network/inventec-mac-config/inventec-mac-config.service
[Unit]
Before=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/inventec-mac-config.sh

[Install]
WantedBy=multi-user.target
