# meta-inventec/meta-transformers-nuv/recipes-phosphor/gpio/gpio-event-handler/smiteimout-handler@.service
[Unit]
Description=SMI Timeout eventhandler

[Service]
ExecStart=/usr/bin/smitimeout %i

[Install]
WantedBy=multi-user.target
