# meta-inventec/meta-transformers-nuv/recipes-phosphor/init/transformers-nuv-init/transformers-nuv-pcie-scanner.service
[Unit]
Description=PCIe Device Scanner for transformers-nuv
After=transformers-nuv-init.service

[Service]
ExecStart=/usr/sbin/transformers-nuv-pcie-scanner.sh
Restart=always

[Install]
WantedBy=multi-user.target
