# meta-inventec/meta-common/recipes-phosphor/flash/configure-firmware-status/service_files/configure-firmware-status.service
[Unit]
Description=configure-firmware-status

[Service]
ExecStart=/usr/bin/configure-firmware-status
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
