# meta-inventec/meta-transformers/recipes-google/acpi-power-state/acpi-power-state-daemon/transformers-check-host-state.service
[Unit]
Description = Check Host State to set fan failsafe speed

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/transformers-check-host-state.sh

[Install]
WantedBy=multi-user.target