# meta-inventec/meta-common/recipes-inventec/util/inventec-util/eeprom-manager.service
[Unit]
Description=Service for accessing EEPROM configurations

[Service]
ExecStart=/usr/bin/eeprom_util

[Install]
WantedBy=multi-user.target
