# meta-inventec/meta-starscream/recipes-phosphor/gpio/gpio-event-handler/gpio-event-handler@.service
[Unit]
Description=GPIO event handle example

[Service]
ExecStart=/usr/bin/gpio-event-handler %i

[Install]
WantedBy=multi-user.target
