# meta-inventec/meta-starscream/recipes-google/acpi-power-state/acpi-power-state-daemon/starscream-host-s0-set-failsafe.service
[Unit]
PartOf=host-s0-state.target
Wants=mapper-wait@-xyz-openbmc_project-settings-fanctrl.service
After=mapper-wait@-xyz-openbmc_project-settings-fanctrl.service

[Service]
Type=oneshot
# let pid control fans, give speed arg as 0
ExecStart=/bin/bash /usr/bin/starscream-states-control.sh 0

[Install]
WantedBy=host-s0-state.target
