# meta-inventec/meta-transformers/recipes-phosphor/gpio/gpio-event-handler/smitimeout-handler@.service
[Unit]
Description=SMI Timeout eventhandler

[Service]
ExecStart=/usr/bin/smitimeout.sh %i

[Install]
WantedBy=multi-user.target
