# meta-inventec/meta-transformers-nuv/recipes-x86/chassis/x86-power-control/host-power-on.service
[Unit]
Description=System unit to for host power on

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/host-power-on.sh

[Install]
WantedBy=host-power-on.target
