# meta-inventec/meta-starscream/recipes-phosphor/ipmi/phosphor-ipmi-flash/verify-image.service
[Unit]
Description=Verify BMC Image File

[Service]
Type=oneshot
ExecStart=/usr/bin/env verify_image.sh

[Install]
RequiredBy=phosphor-ipmi-flash-bmc-verify.target