# meta-inventec/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-psu-update@.service
[Unit]
Description=FLASH PSU %I to PSU

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/dps2400eb-psu-upgrade %i
StandardOutput=journal+console
