# meta-inventec/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-cpld-update@.service
[Unit]
Description=FLASH CPLD %I to CPLD

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/cpldupdate-i2c %i
StandardOutput=journal+console
