# meta-inventec/meta-common/recipes-x86/chassis/x86-power-control/host-power-on.target
[Unit]
Description=System unit to for host power on