# meta-inventec/meta-common/recipes-x86/chassis/x86-power-control/host-power-off.service
[Unit]
Description=System unit to for host power off

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=echo Please add custom command for host power off

[Install]
WantedBy=host-power-off.target
